# diary.service
[Unit]
Description=Diary Web App
Documentation=https://github.com/user/diary_web
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/diary_web
ExecStart=/home/pi/.venv/bin/python main.py
Restart=on-failure
RestartSec=10

# 保护 – 禁止 fork 提权
NoNewPrivileges=yes
PrivateTmp=no
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi/diary_web/data /home/pi/diary_web/config /tmp
ReadOnlyPaths=/home/pi/diary_web/static /home/pi/diary_web/app

# SD 卡保护：减少日志写入
StandardOutput=null
StandardError=journal

# Pi 优化环境变量
Environment=DIARY_PORT=9000
Environment=DIARY_PBKDF2_ITERATIONS=150000

[Install]
WantedBy=multi-user.target
